# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true

[program:php-fpm]
command = /usr/sbin/php-fpm7.3 -c /etc/php
autostart=true
